# systemd/nous-cli.service
[Unit]
Description=Nous CLI Daemon
After=default.target

[Service]
ExecStart=%h/.cargo/bin/nous serve --transport http
Restart=on-failure
RestartSec=5
Environment=RUST_LOG=info
StandardOutput=journal
StandardError=journal
SyslogIdentifier=nous-cli

[Install]
WantedBy=default.target
